# playbooks/legacy/modules/fedora/systemd/kworker-fix/kworker-fix.service
[Unit]
Description=kworker fix

[Service]
Type=oneshot
ExecStart=/usr/bin/kworker-fix

[Install]
WantedBy=multi-user.target
