# modules/backup/files/backup-code.timer
[Unit]
Description=Run backup on code once a day

[Timer]
OnCalendar=daily
Persistent=true 

[Install]
WantedBy=timers.target
